# deploy/smsbot.service
[Unit]
Description=NumberHub Telegram bot (@TheNumberHubBot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/smsbot
ExecStart=/opt/smsbot/.venv/bin/python bot.py
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
# Single-instance: systemd keeps exactly one bot polling the token.
KillMode=control-group

[Install]
WantedBy=multi-user.target
